# roles/jitsi_meet_turnserver/files/override.service
# Ansible managed

[Service]
RestrictAddressFamilies=
SystemCallArchitectures=
SystemCallFilter=
RestartSec=15s
Restart=on-failure
