# roles/jitsi_meet_systemd/files/restart-on-failure.service
# Ansible managed

[Service]
Restart=on-failure
RestartSec=2
