# roles/jitsi_meet_jibri/files/override.service
# Ansible managed

[Service]
ExecStartPre=
ExecStartPre=sudo -u jibri /usr/bin/chromium --timeout=1000 --headless about:blank
